# docker-etc-hosts.service
[Unit]
Description=Monitor Docker events and update /etc/hosts
After=docker.service
Requires=docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/bin/bash -c 'docker events --filter event=start --format "{{.Type}}" | while read event; do /usr/local/bin/docker_etc_hosts; done'
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
